# files/system/etc/skel/.config/autostart/oasis-bottles-init.desktop
[Desktop Entry]
Type=Application
Name=Oasis Bottles Initializer
Comment=Automatically configures the default WINE bottle on first boot
Exec=/usr/libexec/oasis-bottles-init
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true

